# sh/bluetooth/serviceFiles/bluetoothInterpreter.service
[Unit]
Description=RFCOM service
After=bluetoothCommunication.service
Requires=bluetoothCommunication.service

[Service]
ExecStart=sudo sh /usr/bin/bluetoothInterpreter.sh

[Install]
WantedBy=multi-user.target